# services/aprcoind/finish
#!/usr/bin/env sh

aprcoin-cli stop

echo >&2 "APRCoinD exited. code=${1}"

# terminate other services to exit from the container
exec s6-svscanctl -t /var/run/s6/services
